# semantic-visualizer.service
[Unit]
Description=Semantic Distance Visualizer Web Application
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/semantic_distance_visualizer
Environment="PATH=/home/ubuntu/semantic_distance_visualizer/.venv/bin"
ExecStart=/home/ubuntu/semantic_distance_visualizer/.venv/bin/python /home/ubuntu/semantic_distance_visualizer/app.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

